# scripts/deb/postinst
#!/usr/bin/env bash
set -euo pipefail

SHARE=/usr/share/claude-docker
REAL_USER="${SUDO_USER:-}"

case "$1" in
  configure)
    # ── 1. Install Docker if not present ─────────────────────────────────────
    if ! command -v docker >/dev/null 2>&1; then
      echo "==> Installing Docker..."
      apt-get update -qq
      apt-get install -y curl
      curl -fsSL https://get.docker.com | sh
      systemctl enable --now docker || true
    fi

    if [ -n "$REAL_USER" ] && ! id -nG "$REAL_USER" 2>/dev/null | grep -qw docker; then
      usermod -aG docker "$REAL_USER" || true
      echo "    Note: re-login or run 'newgrp docker' for the group change to take effect."
    fi

    # ── 2. Install files to ~/claude-docker (mirrors install.sh) ─────────────
    if [ -n "$REAL_USER" ]; then
      INSTALL_DIR="/home/$REAL_USER/claude-docker"
      LAUNCHER="$INSTALL_DIR/launcher"

      echo "==> Installing claude-docker to $INSTALL_DIR"
      sudo -u "$REAL_USER" mkdir -p \
        "$INSTALL_DIR/home/scratch" \
        "$INSTALL_DIR/certs" \
        "$INSTALL_DIR/comfyui/models" \
        "$INSTALL_DIR/comfyui/output" \
        "$INSTALL_DIR/comfyui/user/default/workflows"
      chmod 700 "$INSTALL_DIR"

      install -o "$REAL_USER" -m 755 "$SHARE/launcher"           "$INSTALL_DIR/launcher"
      install -o "$REAL_USER" -m 644 "$SHARE/docker-compose.yml" "$INSTALL_DIR/docker-compose.yml"
      install -o "$REAL_USER" -m 644 "$SHARE/.env.example"       "$INSTALL_DIR/.env.example"

      if [ ! -f "$INSTALL_DIR/.env" ]; then
        install -o "$REAL_USER" -m 600 "$SHARE/.env.example" "$INSTALL_DIR/.env"
        echo "    Created $INSTALL_DIR/.env from .env.example"
      else
        echo "    Keeping existing $INSTALL_DIR/.env"
      fi

      # ── 3. Add launcher to PATH ───────────────────────────────────────────
      rm -f /tmp/claude-docker-chosen-cmd
      if [ -e /dev/tty ]; then
        sudo -u "$REAL_USER" "$LAUNCHER" --install </dev/tty
      else
        echo "    No interactive terminal — run '$LAUNCHER --install' to add a PATH command"
      fi
      chosen_cmd="$(cat /tmp/claude-docker-chosen-cmd 2>/dev/null || true)"
      rm -f /tmp/claude-docker-chosen-cmd

      cat <<EOF

✓ Installation complete

Next step — edit $INSTALL_DIR/.env:
  • PROJECT_DIRS             — comma-separated host folders, each mounted as /workspace/<name>
                               e.g. PROJECT_DIRS=/home/user/work,/home/user/personal
  • CLAUDE_CODE_OAUTH_TOKEN  — from \`claude setup-token\`, or leave empty to /login inside

Then run: ${chosen_cmd:-$LAUNCHER}

EOF
    fi
    ;;
esac
